# package/postinst
#!/bin/bash

chmod 664 /etc/systemd/system/camerastreamer.service
systemctl enable camerastreamer.service
systemctl start camerastreamer.service
